# super.conf
[supervisord]


[program:app]
command=/root/go/src/backend/backend
directory=/root/go/src/backend/
user=root
autostart=true
